# dot_config/systemd/user/chezmoi-update.service
[Unit]
Description=chezmoi update (pull + apply dotfiles)

[Service]
Type=oneshot
ExecStart=/usr/bin/chezmoi update
